# Packaging/cz.szypowi.brightsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>cz.szypowi.brightsync</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/brightsync</string>
	<key>RunAtLoad</key>
	<true/>
	<!-- Restart on crash only; unconditional KeepAlive makes launchd log an
	     "inherently inefficient" complaint on every registration. -->
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
</dict>
</plist>
